# server/etc/init/301tl.conf
description "301.tl server"
author      "naholyr"
 
start on startup
stop on shutdown
 
script
    # We found $HOME is needed. Without it, we ran into problems
    export HOME="/path/to/home/user"
 
    exec sudo -u user /usr/local/bin/node /path/to/app.js 2>&1 >> /path/to/app.log
end script
